# unsupported/nbssetclk/root.crontab
30 11,2,5, * * * /etc/nbs_refresh > /dev/null
